# supervisord.conf
[program:my_server]
command=/bin/python /home/ec2-user/Cloud7Students/main.py
autostart=true
autorestart=true
stderr_logfile=/var/log/my_server.err.log
stdout_logfile=/var/log/my_server.out.log
