# supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log ; main log file
logfile_maxbytes=50MB               ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10                  ; main logfile backup count; default 10
loglevel=info                       ; log level; default info
pidfile=/var/run/supervisord.pid     ; supervisord pidfile; default supervisord.pid
nodaemon=false                      ; start in foreground if true; default false
minfds=1024                         ; min. avail startup file descriptors; default 1024
minprocs=200                        ; min. avail process descriptors; default 200
user=root                            ; setuid to this UNIX account to run supervisord
childlogdir=/var/log/supervisord     ; 'AUTO' child log dir, default $TEMP
strip_ansi=false  

[program:my_server]
command=/bin/python /home/ec2-user/Cloud7Students/main.py
autostart=true
autorestart=true
stderr_logfile=/var/log/my_server.err.log
stdout_logfile=/var/log/my_server.out.log
